# rainshine.service
[Unit]
Description=Rainshine DMX Shader
After=network.target olad.service
Wants=olad.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/rainshine-env/bin/python3 /home/pi/rainshine_dmx.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
